# systemd_examplefiles/speedtest.target
[Unit]
Description=Speedtest service
StopWhenUnneeded=yes